# deploy/kfd-quote-mailer.service
[Unit]
Description=KFD quote email delivery service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kfd-mailer
Group=kfd-mailer
WorkingDirectory=/opt/kfd-quote-mailer
EnvironmentFile=/etc/kfd-quote-mailer.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
